# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=www-data
logfile=/var/www/storage/tmp/supervisord.log
pidfile=/var/www/storage/tmp/supervisord.pid

[program:php-fpm]
command=php-fpm
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:schedule]
command=php /var/www/docker/schedule/schedule.php "php /var/www/artisan schedule:run > /var/www/storage/tmp/schedule.log 2>&1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0